# runit/log/run
#!/bin/bash
set -e pipefail

exec runit-log dyrwood
